# rootfs/etc/s6-overlay/s6-rc.d/fcgiwrap/up
if { spawn-fcgi -s /var/run/fcgiwrap/cgit.sock -u git -g git -- /usr/bin/fcgiwrap }
chown nginx:nginx /var/run/fcgiwrap/cgit.sock
