# rootfs/etc/s6-overlay/s6-rc.d/git-daemon/run
#!/command/execlineb -P
/usr/bin/git daemon --verbose --reuseaddr --user=git --base-path=/var/lib/git/repositories /var/lib/git/repositories
