# rootfs/etc/s6-overlay/s6-rc.d/git-daemon/run
#!/command/execlineb -P
s6-setuidgid git /usr/bin/git daemon --reuseaddr --base-path=/data /data
